# scripts/bedriftsgrafen-regnskap-sync.service
[Unit]
Description=Bedriftsgrafen Continuous Regnskap Sync Service
Documentation=https://github.com/bedriftsgrafen/bedriftsgrafen.no
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

# Don't start until Docker is fully ready
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=k1sso
Group=k1sso
WorkingDirectory=/home/k1sso/bedriftsgrafen.no

# Use docker compose exec for proper container lifecycle management
ExecStart=/usr/bin/docker compose -f /home/k1sso/bedriftsgrafen.no/docker-compose.prod.yml exec -T backend python3 /app/ops_scripts/continuous_regnskap_sync.py

# Restart policy
Restart=always
RestartSec=60

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bedriftsgrafen-regnskap-sync

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/k1sso/bedriftsgrafen.no/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
